# files/rpi/etc/uci-defaults/97-network
#!/bin/sh

uci -q batch << EOF
set network.lan.proto='dhcp'
set dhcp.lan.ignore='1'
EOF
